# 5semester/SystemProgramming/linux/lab9/Service09d.service
[Unit]
Description=watch directory daemon

[Service]
Type=exec
ExecStart=/usr/sbin/Service09d

[Install]
WantedBy=multi-user.target

